# docker/dev/php-cli/supervisor.d/crm-send-reminders.ini
[program:crm-send-reminders]
command=php /var/www/bin/app.php crm:send-reminders
directory=/var/www
autostart=false
autorestart=false
startsecs=0
startretries=0
stdout_logfile=/var/log/supervisor/crm-send-reminders.log
stderr_logfile=/var/log/supervisor/crm-send-reminders-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
