# docker/dev/php-cli/supervisor.d/cron.ini
[program:cron]
command=crond -f -l 2
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/cron.log
stderr_logfile=/var/log/supervisor/cron-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
