# docker/dev/php-cli/supervisor.d/notify-subscription-expired.ini
[program:notify-subscription-expired]
command=php /var/www/bin/app.php notify-subscription-expired
directory=/var/www
autostart=false
autorestart=false
startsecs=0
startretries=0
stdout_logfile=/var/log/supervisor/notify-subscription-expired.log
stderr_logfile=/var/log/supervisor/notify-subscription-expired-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
