# docker/dev/php-cli/supervisor.d/notify-subscription-expiring.ini
[program:notify-subscription-expiring]
command=php /var/www/bin/app.php notify-subscription-expiring
directory=/var/www
autostart=false
autorestart=false
startsecs=0
startretries=0
stdout_logfile=/var/log/supervisor/notify-subscription-expiring.log
stderr_logfile=/var/log/supervisor/notify-subscription-expiring-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
