# docker/dev/php-cli/supervisor.d/parse-metro-stations.ini
[program:parse-metro-stations]
command=php /var/www/bin/app.php parse-metro-stations
directory=/var/www
autostart=true
autorestart=false
startsecs=0
startretries=0
stdout_logfile=/var/log/supervisor/parse-metro-stations.log
stderr_logfile=/var/log/supervisor/parse-metro-stations-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
