# docker/dev/php-cli/supervisor.d/parse-yandex.ini
[program:parse-yandex]
command=php /var/www/bin/app.php parse-yandex
directory=/var/www
autostart=false
autorestart=true
startsecs=10
startretries=20
stdout_logfile=/var/log/supervisor/parse-yandex.log
stderr_logfile=/var/log/supervisor/parse-yandex-error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=10MB
