# docker/dev/php-cli/supervisor.d/photo-tasks.ini
[program:photo-tasks]
command=php /var/www/bin/app.php photo-tasks
directory=/var/www
autostart=true
autorestart=true
startsecs=10
startretries=20
stdout_logfile=/var/log/supervisor/photo-tasks.log
stderr_logfile=/var/log/supervisor/photo-tasks-error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=10MB
