# docker/dev/php-cli/supervisor.d/transfer-listings.ini
[program:transfer-listings]
command=php /var/www/bin/app.php transfer-listings
directory=/var/www
autostart=true
autorestart=true
startsecs=10
startretries=20
stdout_logfile=/var/log/supervisor/transfer-listings.log
stderr_logfile=/var/log/supervisor/transfer-listings-error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=10MB
