# docker/dev/php-cli/supervisor.d/update-expired-subscriptions.ini
[program:update-expired-subscriptions]
command=php /var/www/bin/app.php update-expired-subscriptions
directory=/var/www
autostart=false
autorestart=false
startsecs=0
startretries=0
stdout_logfile=/var/log/supervisor/update-expired-subscriptions.log
stderr_logfile=/var/log/supervisor/update-expired-subscriptions-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
